# garagedoor-mqtt.service
[Unit]
After=networking.service

[Service]
ExecStart=/usr/bin/python3 /home/pi/garagedoor-mqtt/garagedoor.py

[Install]
WantedBy=default.target
